# debian/polaric-webapp.cron.hourly
#!/bin/bash

# cleanup of map files
# This should be moved to the polaric-webapp package! 

/usr/bin/find /var/www/aprs/kacache2/. -size 1714c -exec rm {} \;
/usr/bin/find /var/www/aprs/kacache2/. -size 1715c -exec rm {} \;
/usr/bin/find /var/www/aprs/kacache2/. -size 1716c -exec rm {} \;
/usr/bin/find /var/www/aprs/kacache2/. -size 1717c -exec rm {} \;
chown -R www-data.www-data /var/www/aprs/kacache2/